# integration-framework/moirai-infrastructure/ansible/roles/flask-s3-zip-service/app/s3zip.service
[Unit]
Description=Flask Service

[Service]
WorkingDirectory=/home/ubuntu/
ExecStart=/home/ubuntu/s3zip.py
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target