# integration-framework/moirai-infrastructure/ansible/roles/flask-githook-service/app/flask.service
[Unit]
Description=Flask Service

[Service]
WorkingDirectory=/home/ubuntu/
ExecStart=/home/ubuntu/git-hook.py
Restart=always
User=ubuntu

[Install]
WantedBy=multi-user.target